# deploy/systemd/ultra-codeexec.service
[Unit]
Description=Ultra code execution service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/ultra/codeexec-service.env
ExecStartPre=/bin/mkdir -p ${CODEEXEC_JOB_ROOT:-/srv/ultra/codeexec-service/jobs} ${CODEEXEC_ARTIFACT_ROOT:-/srv/ultra/codeexec-service/artifacts}
ExecStartPre=-/usr/bin/docker rm -f ultra-codeexec-service
ExecStart=/bin/sh -lc '/usr/bin/docker run --rm --name ultra-codeexec-service --env-file /etc/ultra/codeexec-service.env -p ${CODEEXEC_SERVICE_BIND_ADDRESS:-0.0.0.0}:${CODEEXEC_SERVICE_PORT:-8020}:8020 -v /usr/bin/docker:/usr/bin/docker:ro -v /var/run/docker.sock:/var/run/docker.sock -v ${CODEEXEC_JOB_ROOT:-/srv/ultra/codeexec-service/jobs}:${CODEEXEC_JOB_ROOT:-/srv/ultra/codeexec-service/jobs} -v ${CODEEXEC_ARTIFACT_ROOT:-/srv/ultra/codeexec-service/artifacts}:${CODEEXEC_ARTIFACT_ROOT:-/srv/ultra/codeexec-service/artifacts} ultra-codeexec-service:current'
ExecStop=/usr/bin/docker stop ultra-codeexec-service
Restart=always
RestartSec=5
TimeoutStartSec=0
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
